# tools/docker/soap/supervisor/conf.d/soap.conf
[program:soap]
command=java api.java.openiss.ws.soap.endpoint.ServicePublisher
directory=/usr/src/app
autorestart=true